# hass-desktop-sensor.service
[Unit]
Description=Home Assistant Desktop Sensor

[Service]
Type=simple
ExecStart=/path/to/hass-desktop-sensor

[Install]
WantedBy=multi-user.target
